# nfsroot/var/lib/dpkg/info/fai-nfsroot.conffiles
/etc/apt/apt.conf.d/90fai
/etc/init.d/fai-abort
/etc/init/fai.conf
