# nfsroot/var/lib/dpkg/info/fai-nfsroot.list
/.
/etc
/etc/init.d
/etc/init.d/fai-abort
/etc/fai
/etc/apt
/etc/apt/apt.conf.d
/etc/apt/apt.conf.d/90fai
/etc/init
/etc/init/fai.conf
/usr
/usr/lib
/usr/lib/dracut
/usr/lib/dracut/modules.d
/usr/lib/dracut/modules.d/80fai-autodiscover
/usr/lib/dracut/modules.d/80fai-autodiscover/fai-scan.sh
/usr/lib/dracut/modules.d/80fai-autodiscover/module-setup.sh
/usr/lib/dracut/modules.d/80fai-autodiscover/parse-fai.sh
/usr/lib/fai
/usr/lib/fai/create_resolv_conf
/usr/lib/fai/dhclient-perl
/usr/lib/fai/fai-mount-disk
/usr/lib/fai/task_sysinfo
/usr/lib/fai/load_keymap_consolechars
/usr/lib/fai/mount2dir
/usr/share
/usr/share/fai
/usr/share/fai/dhclient-fai-script
/usr/share/fai/dhclient-fai.conf
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/fai-nfsroot
/usr/share/doc
/usr/share/doc/fai-nfsroot
/usr/share/doc/fai-nfsroot/changelog.gz
/usr/share/doc/fai-nfsroot/NEWS.Debian.gz
/usr/share/doc/fai-nfsroot/copyright
/usr/sbin
/usr/sbin/faireboot
/usr/bin
/usr/bin/policy-rc.d.fai
