# nfsroot/var/lib/dpkg/info/smartmontools.conffiles
/etc/default/smartmontools
/etc/init.d/smartmontools
/etc/smartd.conf
/etc/smartmontools/run.d/10mail
/etc/smartmontools/run.d/10powersave-notify
