# nfsroot/var/lib/dpkg/info/smartmontools.prerm
#!/bin/sh
set -e
# Automatically added by dh_installdeb
dpkg-maintscript-helper rm_conffile /etc/smartd_warning.sh 6.3+svn4002-2~ -- "$@"
# End automatically added section
# Automatically added by dh_installinit
if [ -x "/etc/init.d/smartmontools" ]; then
	invoke-rc.d smartmontools stop || exit $?
fi
# End automatically added section
