# nfsroot/var/lib/dpkg/info/smartmontools.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/smartd.service
/var
/var/lib
/var/lib/smartmontools
/var/lib/smartmontools/drivedb
/var/lib/smartmontools/drivedb/drivedb.h
/usr
/usr/sbin
/usr/sbin/smartctl
/usr/sbin/smartd
/usr/sbin/update-smart-drivedb
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/smartmontools
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/update-smart-drivedb.8.gz
/usr/share/man/man8/smartd.8.gz
/usr/share/man/man8/smartctl.8.gz
/usr/share/man/man5
/usr/share/man/man5/smartd.conf.5.gz
/usr/share/doc
/usr/share/doc/smartmontools
/usr/share/doc/smartmontools/ChangeLog.gz
/usr/share/doc/smartmontools/ChangeLog-5.0-6.0.gz
/usr/share/doc/smartmontools/COPYING.gz
/usr/share/doc/smartmontools/INSTALL.gz
/usr/share/doc/smartmontools/NEWS.gz
/usr/share/doc/smartmontools/smartd.conf.gz
/usr/share/doc/smartmontools/changelog.Debian.gz
/usr/share/doc/smartmontools/changelog.Debian.amd64.gz
/usr/share/doc/smartmontools/NEWS.Debian.gz
/usr/share/doc/smartmontools/changelog.gz
/usr/share/doc/smartmontools/copyright
/usr/share/doc/smartmontools/README.Debian
/usr/share/doc/smartmontools/badblockhowto.html
/usr/share/doc/smartmontools/examples
/usr/share/doc/smartmontools/examples/Example6
/usr/share/doc/smartmontools/examples/Example5
/usr/share/doc/smartmontools/examples/Example4
/usr/share/doc/smartmontools/examples/Example3
/usr/share/doc/smartmontools/examples/Example2
/usr/share/doc/smartmontools/examples/Example1
/usr/share/doc/smartmontools/examples/README
/usr/share/doc/smartmontools/WARNINGS
/usr/share/doc/smartmontools/TODO
/usr/share/doc/smartmontools/README
/usr/share/doc/smartmontools/AUTHORS
/usr/share/smartmontools
/usr/share/smartmontools/smartd-runner
/usr/share/smartmontools/smartd_warning.sh
/etc
/etc/init.d
/etc/init.d/smartmontools
/etc/default
/etc/default/smartmontools
/etc/smartd.conf
/etc/smartmontools
/etc/smartmontools/smartd_warning.d
/etc/smartmontools/run.d
/etc/smartmontools/run.d/10powersave-notify
/etc/smartmontools/run.d/10mail
/lib/systemd/system/smartmontools.service
